# roles/wireguard/files/update_duckdns.service
[Unit]
Description=Update IP in DuckDNS

[Service]
Type=oneshot
ExecStart=/usr/bin/update_duckdns.sh

[Install]
WantedBy=multi-user.target